# audio/setup/templates/soundbar-keepalive.service.tmpl
[Unit]
Description=Soundbar keepalive tone (18kHz)
After=pipewire-pulse.service
Requires=pipewire-pulse.service

[Service]
ExecStart=ffmpeg \
  -f lavfi -i "sine=frequency=18000:sample_rate=44100" \
  -af "volume=0.05" \
  -f pulse default
Environment=PULSE_SINK={{.PrimarySink}}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
